# configs/supervisor.ddserver.conf
[program:ddserver]
command=/usr/bin/ddserver-bundle
redirect_stderr=false
