# configs/supervisor.pdns.conf
[program:pdns]
command=/usr/sbin/pdns_server
